# server/crontab/autoMail.sh
#! /bin/sh
curl http://test.f2e.taobao.net/server/API/doMailList.php\?action=mail\&token=9527